# systemd/claude-code-router.service
[Unit]
Description=claude-code-router (Anthropic <-> local model routing proxy)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/claude-code-router/proxy.py
Restart=on-failure
RestartSec=2

# Run as an unprivileged dynamic user; loopback-only listener so no caps needed.
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Defaults; override with `systemctl edit` or an EnvironmentFile if needed.
Environment=CCR_LISTEN_HOST=127.0.0.1
Environment=CCR_LISTEN_PORT=8082
Environment=CCR_LOCAL_HOST=127.0.0.1
Environment=CCR_LOCAL_PORT=8080

# Optional: rewrite the model name sent to the local backend.
# Useful when the client sends a claude-* name but your local server
# only knows its own model identifier.
#Environment=CCR_LOCAL_MODEL=my-local-model-name

# Optional: comma-separated substrings — any claude-* model whose name
# contains one of these is routed to the local backend instead of Anthropic.
# Example: "haiku" routes claude-haiku-* locally, sonnet/opus go to Anthropic.
# Leave unset to route all claude-* models to Anthropic (legacy behaviour).
#Environment=CCR_LOCAL_MODELS=haiku

[Install]
WantedBy=multi-user.target
